# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
set -e

umask ${UMASK:-077}
exec nginx -g "daemon off;"
